# etc/profile.d/chef.zsh
#!zsh

# Use vagrant-digitalocean
export SSL_CERT_FILE=/usr/local/opt/curl-ca-bundle/share/ca-bundle.crt

# Set up knife aliases
alias ku='knife upload'
alias bki='berks install'
alias bku='berks upload'
alias bkc='berks cookbook'
