# etc/profile.d/chef.zsh
#!zsh

# Set up knife aliases
alias ku='knife upload'
alias bki='berks install'
alias bku='berks upload'
alias bkc='berks cookbook'
