# etc/profile.d/config.zsh
#
# Basic ZSH shell configuration
#


# Use 'viins' mappings for *most* inline command editing
#bindkey -v
