# etc/profile.d/config.zsh
#
# Basic ZSH shell configuration
#

# Edit files with Vim
export EDITOR='vim'

# View files with Less
export PAGER='less -r'
