# etc/profile.d/ember.zsh
alias en='ember new'
alias ena='ember addon'
alias ea=ena
alias eg='ember generate'
alias es='ember server'
alias ein='ember install:npm'
alias eib='ember install:bower'
alias eia='ember install:addon'
