# etc/profile.d/git.zsh
#
# Low-character aliases for working with Git.
#

alias s='git status'
alias gl='git blog'
alias ga='git add'
alias gc='git commit'
alias gr='git rebase'
alias gp='git push'
alias gu='git pull --rebase'
alias gm='git merge'
alias gur=gu
alias gd='git diff'
alias gmt='git mergetool'
alias gra='git rebase --abort'
alias grc='git rebase --continue'
alias grs='git rebase --skip'
alias garc='ga . && grc'
alias gg='git grep'
alias gco='git checkout'
alias gth='git reset --hard'
alias gthh='gth HEAD'
alias gb='git branch'
