# etc/profile.d/make.zsh
#!zsh

alias m='make'
alias mt='make test'
alias mw='make watch'
alias mi='make install'
alias mb='make build'
alias cmmi='./configure && m && mi'
