# etc/profile.d/plugins.zsh
# Use Antigen to manage plugins.

source /usr/local/share/homer/plugins/antigen.zsh

antigen bundle zsh-users/zsh-syntax-highlighting
antigen bundle kennethreitz/autoenv

