# etc/profile.d/powerline.zsh
#!/bin/zsh
#
# Use powerline for the prompt

source /usr/local/lib/python3.4/site-packages/powerline/bindings/zsh/powerline.zsh
