# etc/profile.d/rails.zsh
#!zsh
#
# Helpful Rails aliases

alias rg='bin/rails generate'
alias rdb='bin/rails db'
alias rc='bin/rails console'
alias rt='bin/rake test'
alias ts='bin/rspec'
alias rdm='bin/rake db:migrate'

# pow
alias powder=powify
alias rs='powify restart'
alias rl='powify logs'
