# etc/profile.d/ssh.zsh
# SSH agent forwarding
zstyle :omz:plugins:ssh-agent agent-forwarding on
