# etc/profile.d/travis.zsh
#!/bin/zsh
#
# Configure the Travis gem.

[ -f /Users/necromancer/.travis/travis.sh ] && source /Users/necromancer/.travis/travis.sh
